# zigstar-cc2652-flasher/rootfs/etc/s6-overlay/s6-rc.d/cc2652-flasher/up
/etc/s6-overlay/scripts/cc2652-flasher-up
